# init/creekd.service
[Unit]
Description=creekd — multi-tenant process supervisor
Documentation=https://github.com/solcreek/creekd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/creekd
Restart=on-failure
RestartSec=5

# Identity. Stable UID required so /var/lib/creekd ownership
# survives daemon restarts and package upgrades. DynamicUser=no per
# DESIGN-self-host-state.md §"creekd privilege model & systemd
# hardening".
User=creekd
Group=creekd
DynamicUser=no

# Hardening — every directive below is required (not advisory).
# `creek host doctor` (lands with #23 CLI) reports
# `systemd_hardening_drift` if any are missing or weakened on the
# running host. Validator implementation: internal/hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service ~@privileged ~@resources

# Capabilities. creekd binds 80/443 for ACME without full root;
# CAP_NET_BIND_SERVICE is the minimum that allows that.
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# No core dumps — protects in-memory secrets (#7 hostkey, future
# KEK envelope material) from landing on disk on crash.
LimitCORE=0
LimitNOFILE=65536

# Writable paths. ProtectSystem=strict makes the rest of the FS
# read-only; these two override.
ReadWritePaths=/var/lib/creekd /var/log/creekd

# State + logs
Environment=CREEKD_STATE_DIR=/var/lib/creekd/state
Environment=CREEKD_LOG_DIR=/var/lib/creekd/logs

# Admin API (loopback only by default)
Environment=CREEKD_ADMIN_ADDR=127.0.0.1:9080

# Dispatch proxy (public-facing, Caddy terminates TLS in front)
Environment=CREEKD_DISPATCH_ADDR=127.0.0.1:9000

# cgroup v2 isolation
Environment=CREEKD_CGROUP_PARENT=creekd.slice
Environment=CREEKD_DEFAULT_MEMORY_HIGH=256M
Environment=CREEKD_DEFAULT_MEMORY_MAX=1G

# Override via /etc/creekd/env if needed
EnvironmentFile=-/etc/creekd/env

[Install]
WantedBy=multi-user.target
